# config-HAPU/etc/rc.local
# Mount
mount /dev/sda3 /data

# Scripts
/root/ipset.sh start
/etc/init.d/telegraf start
/root/throttle.sh start
/etc/init.d/web-backend start

# End programs
exit 0

